# openwrt/testmesh/files/etc/hotplug.d/iface/30-testmesh
#!/bin/sh

. /lib/functions.sh

handle_ifname_list() {
	local ifname="$1" cfg="$2" control="/tmp/testmesh-ctl"

	config_get control $cfg "control"

	if [ "$ifname" = "$INTERFACE" ]; then
		if [ "$ACTION" = "ifup" ]; then
			testmesh-ctl -s $control interface-add "$DEVICE"
		fi
		if [ "$ACTION" = "ifdown" ]; then
			testmesh-ctl -s $control interface-del "$DEVICE"
		fi
	fi
}

handle_testmesh_section() {
	config_list_foreach "$1" ifname handle_ifname_list "$1"
}

config_load testmesh
config_foreach handle_testmesh_section testmesh
